# telegram-bot/lato-telegram-bot.service
[Unit]
Description=Lato Telegram Bot (departman input isleyici, Claude Sonnet 5)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/lato-knowledge/telegram-bot
EnvironmentFile=/root/.hermes/profiles/lato/.env
# claude CLI'nin bulunmasi + ~/.claude abonelik kimligi icin
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/lib/node_modules/.bin
ExecStart=/usr/bin/env python3 lato_telegram_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
